# patches/system/console-blank.patch/build/files/etc/local.d/console-blank.start
# put consoles into powersave just like we do in Xorg
# NOTE: this is also set on the kernel commandline, consoleblank=300
setterm --blank 5

